# deploy/com.kalshibot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Kalshi Bot — macOS LaunchAgent

  Runs the weather strategy daily at 7 AM CT (= 8 AM ET).
  Stops after 9 hours, covering the full trading day through ~4 PM CT.

  Why 7 AM CT / 8 AM ET?
    - Next-day markets open the evening before; the morning NWS model run
      (released ~5-6 AM ET) contains overnight forecast updates. Starting at
      8 AM ET means we catch markets before the rest of the market reacts to
      the new NWS data.
    - σ tightens at 11 AM ET (early = σ=3.0°F, late = σ=2.0°F), so morning
      entries get wider signals but also earlier pricing inefficiencies.
    - 9-hour window ensures same-day markets that close by 4-5 PM ET are covered.

  INSTALL:
      cp deploy/com.kalshibot.plist ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.kalshibot.plist

  CHECK STATUS:
      launchctl list | grep kalshibot

  UNINSTALL:
      launchctl unload ~/Library/LaunchAgents/com.kalshibot.plist
      rm ~/Library/LaunchAgents/com.kalshibot.plist

  MANUAL RUN (test first — dry run):
      python deploy/run_bot.py --strategy weather_edge
  MANUAL RUN (live, one scan):
      python deploy/run_bot.py --strategy weather_edge --live --once
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kalshibot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/alexgarrison/KalshiBotClaude/deploy/run_bot.py</string>
        <string>--strategy</string>
        <string>weather_edge</string>
        <string>--live</string>
        <string>--hours</string>
        <string>9</string>
    </array>

    <!-- 7 AM CT = 8 AM ET (CDT, UTC-5). Adjust Hour if you change timezones. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/alexgarrison/KalshiBotClaude</string>

    <key>StandardOutPath</key>
    <string>/Users/alexgarrison/KalshiBotClaude/logs/launchagent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alexgarrison/KalshiBotClaude/logs/launchagent.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alexgarrison</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
